# luci-app-seaweedfs/files/etc/uci-defaults/90-seaweedfs-setup
#!/bin/sh
mkdir -p /opt/seaweedfs
/etc/init.d/seaweedfs enable 2>/dev/null || true
exit 0

